# backend/deploy/api-radar.service
[Unit]
Description=API Radar do Povo (Laravel)
After=network.target

[Service]
Type=simple
User=www
Group=www
WorkingDirectory=/www/wwwroot/radardopovo.com/backend
Environment=APP_ENV=production
Environment=APP_HOST=127.0.0.1
Environment=APP_PORT=8081
ExecStart=/usr/bin/php artisan serve --host=127.0.0.1 --port=8081
Restart=always
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
